# ops/deployment/zaki-api.service
[Unit]
Description=ZAKI AI CFO Backend (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/Mumtaz/apps/zaki/backend
Environment=PYTHONPATH=/opt/Mumtaz/apps/zaki/backend
EnvironmentFile=/opt/zaki-server/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zaki-api

[Install]
WantedBy=multi-user.target
