# c2p-delivery-system/deploy/delivery-api.service.tmpl
[Unit]
Description=C2P delivery-api
After=network.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__API_DIR__
EnvironmentFile=__API_DIR__/.env
ExecStart=__API_DIR__/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8800
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
