# platform/ops/systemd/mumtaz-platform.service
[Unit]
Description=Mumtaz Control Panel backend (app.mumtaz.digital)
After=network.target postgresql.service
[Service]
WorkingDirectory=/opt/mumtaz/platform/backend
ExecStart=/opt/mumtaz/platform/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8002 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
