# platform/ops/systemd/zaki.service
[Unit]
Description=ZAKI AI CFO backend (Mumtaz)
After=network.target postgresql.service
[Service]
WorkingDirectory=/opt/mumtaz/zaki/backend
ExecStart=/opt/mumtaz/zaki/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
